# dispatch_scripts/nanoclaw-heartbeat.timer
[Unit]
Description=PhilClaw Heartbeat Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=900
AccuracySec=30

[Install]
WantedBy=timers.target
